# files/supervisor_crond.conf

[program:crond]
command = crond -n -s
autorestart = true
autostart = true
stdout_logfile = /dev/stdout
stdout_logfile_maxbytes = 0
stderr_logfile = /dev/stderr
